# launchd/com.football-predictor.odds-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.football-predictor.odds-poll</string>

    <!-- Runs through the backend container (docker compose exec) like every
         other job — see scripts/run_odds_poll.sh. The old form invoked a
         non-existent .venv python directly and failed on every fire. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJ_DIR__/scripts/run_odds_poll.sh</string>
    </array>

    <!-- Run every 3 hours at the top of the hour -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJ_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__PROJ_DIR__</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/odds-poll-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/odds-poll-stderr.log</string>
    <!-- No KeepAlive: this is a scheduled job (StartCalendarInterval). KeepAlive
         with SuccessfulExit=false would relaunch it on every failure → a
         restart storm. The 3-hourly schedule is the only trigger. -->
</dict>
</plist>
