# launchd/com.football-predictor.prematch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.football-predictor.prematch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJ_DIR__/scripts/run_prematch.sh</string>
    </array>

    <!-- Run at 15:00 local time — ~2h before typical evening kick-offs.
         Refreshes today's predictions with closing-line odds (sharper market). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/prematch-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/prematch-stderr.log</string>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
</dict>
</plist>
